# image-extras/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/root/autoprovision-stage1.sh

exit 0
